# debian/grac-device-daemon.prerm
#!/bin/sh
set -e

systemctl stop grac-device-daemon

systemctl disable grac-device-daemon
