# debian/grac-device-daemon.postinst
#!/bin/sh
/bin/systemctl daemon-reload
/bin/systemctl enable grac-device-daemon
/bin/systemctl start grac-device-daemon
/bin/systemctl reload grac-device-daemon
