# debian/preinst
#!/bin/sh
/bin/systemctl stop grac-device-daemon
/bin/systemctl disable grac-device-daemon
/bin/rm -rf /etc/gooroom/grac.d/printers

